# etc/crontab
* * * * * /usr/local/bin/php /app/bin/notify.php > /var/log/notify.log 2>&1
